# planetlab/ipfw.cron
# Runs every 5 minutes and clean ipfw expired rules
# $Id: ipfw.cron 4704 2010-01-11 10:19:56Z marta $
*/5 * * * * root     echo "super killexpired" | /vsys/ipfw-be root > /dev/null 2>&1
